# systemd/scribe.service.in
[Unit]
Description=Scribe hold-to-talk transcription daemon
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=@BINDIR@/scribe
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
